# recipes-lafon/application/files/40-lon.rules
SUBSYSTEM=="tty",ATTRS{idVendor}=="0920", ATTRS{idProduct}=="5550", RUN+="/root/easy/lonifd -d /dev/%k" , OWNER="E@sy"