# recipes-lafon/application/files/magfleet-application.target
[Unit]
Description= all Magfleet applications
Wants= magfleet-app.service magfleet-ihm.service magfleet-ctrl.service lon.service magfleet-jrnl.service

[Install]
WantedBy=multi-user.target
Alias=application.target