# recipes-lafon/application/files/magfleet-mid.target
[Unit]
Description= Magfleet MID applications
Wants= lon.service magfleet-jrnl.service

[Install]
WantedBy=multi-user.target
Alias=mid.target