# recipes-lafon/gestion-carte-sd/files/mount-sd@.service
[Unit]
Description=Mount sd card Service

[Service]
Type=notify
ExecStart=/usr/bin/add_sd_script -p /dev/%I
