# recipes-lafon/gestion-carte-sd/files/umount-sd@.service
[Unit]
Description=Unmount sd card Service

[Service]
Type=notify
ExecStart=/usr/bin/remove_sd_script -p /dev/%I
